# packaging/command-line-assistant.service
[Unit]
Description=Command Line Assistant Daemon
Documentation=https://github.com/rhel-lightspeed/command-line-assistant
After=network.target

[Service]
Type=simple
User=command-line-assistant
Group=command-line-assistant
ExecStart=/usr/bin/clad
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security settings
PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/command-line-assistant

[Install]
WantedBy=multi-user.target

